# apps/agent/deploy/pisignage-agent.service
[Unit]
Description=PiSignage Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pisignage-agent
EnvironmentFile=/etc/pisignage/agent.env
Environment=PISIGNAGE_STATE_DIR=/etc/pisignage
# Wayland session for the Chromium kiosk (adjust UID if not 1000).
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/usr/bin/node /opt/pisignage-agent/dist/index.js
Restart=always
RestartSec=3

[Install]
WantedBy=graphical-session.target
